# vendor/github.com/catalyzeio/paas-orchestration/examples/runit/registry/run
#!/bin/sh -e
exec 2>&1

USER=nobody

exec chpst -u ${USER} /opt/orch/bin/regserver \
    -v=4
